# recording/matrix/supervisor-matrix.conf
[program:matrix_mqtt]
command=node index.js

directory=/home/pi/transfer-learning-playground/recording/matrix
autostart=true
autorestart=true
stopasgroup=true
