# recording/thunderboard/supervisor-thunderboard.conf
[program:matrix_mqtt]
command=sudo python tbsense_scan.py

directory=/home/pi/transfer-learning-playground/recording/thunderboard
autostart=true
autorestart=true
stopasgroup=true
